# etc/services.d/sshd/run
#!/bin/sh

/usr/sbin/sshd -o "LogLevel VERBOSE" -D -E /dev/stdout -f /dev/shm/sshd_config
